# mobile-terminal.service.example
[Unit]
Description=Mobile Terminal Overlay
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/mobile-terminal-overlay
Environment=PATH=/path/to/mobile-terminal-overlay/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/path/to/mobile-terminal-overlay/venv/bin/mobile-terminal --session claude --port 8765
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

# Installation:
# 1. Copy this file to /etc/systemd/system/mobile-terminal.service
# 2. Replace YOUR_USERNAME with your username
# 3. Replace /path/to/mobile-terminal-overlay with the actual path
# 4. Run: sudo systemctl daemon-reload
# 5. Run: sudo systemctl enable mobile-terminal
# 6. Run: sudo systemctl start mobile-terminal
