# mobile-terminal.service
[Unit]
Description=Mobile Terminal Overlay
After=network.target

[Service]
Type=simple
User=gbons
WorkingDirectory=/home/gbons/dev/mobile-terminal-overlay
Environment=PATH=/home/gbons/dev/mobile-terminal-overlay/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/gbons/dev/mobile-terminal-overlay/venv/bin/mobile-terminal --session claude --no-auth --port 8765
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
